# bms-reader.service
[Unit]
Description=BMS Reader for LiFePO4 batteries
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
Group=dialout
WorkingDirectory=/home/pi/bms-reader
Environment=PATH=/home/pi/bms-reader/venv/bin
ExecStart=/home/pi/bms-reader/venv/bin/python /home/pi/bms-reader/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/bms-reader
ProtectHome=true

[Install]
WantedBy=multi-user.target
